# buildroot-external-xaptum/board/xaprw001/rootfs_overlay/etc/systemd/system/xap-usb-gadget.service.d/50-after-macset.conf
[Unit]
After=macset@wlan0.service
After=sys-subsystem-net-devices-wlan0.device
