# buildroot-external-xaptum/board/xaprc001/initramfs_overlay/etc/systemd/system/data.mount
[Unit]
Description=Mount Persistent Data Directory

[Mount]
What=/dev/mmcblk0p4
Where=/data
Type=ext4
Options=defaults

[Install]
WantedBy=local-fs.target
