# buildroot-external-xaptum/board/xaprc-common/rootfs_overlay/etc/systemd/system/secure-host.target
[Unit]
Description=Secure Host Mode
Requires=multi-user.target
Conflicts=passthrough.target secure-lan.target
After=multi-user.target
AllowIsolate=yes
