# buildroot-external-xaptum/board/xaprw001/rootfs_overlay/usr/lib/systemd/system/connman-wifi-config@.service
[Unit]
Description=Configure Connman Wi-Fi credentials for %i
Before=connman.service

[Service]
Type=oneshot
ExecStart=/bin/ln -s -f /data/connman/%i/wifi.config /var/lib/connman/wifi.config
